# config/includes.chroot/etc/rc.local
#!/bin/sh -e 

# Update dhd, if possible
su jessica -c 'cd ~/.dhd ; git pull https://younix.us/cgit/cgit.cgi/dhd.git/' &


exit 0
